# src/update-motd.d/00-header
#!/bin/sh
if [ -f /etc/os-release ]
then
    . /etc/os-release
    echo
    echo "Welcome to $PRETTY_NAME"
fi
