# src/update-motd.d/10-fortune
#!/bin/sh
[ -x /usr/bin/fortune ] && echo && /usr/bin/fortune
