# src/update-motd.d/10-fortune
#!/bin/sh
[ -d /usr/games ] && PATH=/usr/games:$PATH
if type fortune > /dev/null 2>&1
then
	echo
	fortune
fi
