# contrib/archlinux/motd.service
[Unit]
Description=Dynamic MOTD updater service

[Service]
Type=oneshot
ExecStart= /usr/lib/update-motd/update-motd
